# etc/services.d/smbd/run
#!/bin/bash
/usr/sbin/smbd -F -l /config/logs -s /config/smb.conf