# etc/services.d/nmbd/run
#!/bin/bash
/usr/sbin/nmbd -F -l /config/logs -s /config/smb.conf